# config/supervisor.conf
[supervisord]
nodaemon=true

[program:cellect]
user=root
command=/opt/start_puma
redirect_stderr=true
autorestart=true
